# supervisord.conf
[supervisord]
nodaemon=true

[program:languages]
directory=/app/languages
command=/usr/local/bin/node build/index.js

[program:vospel]
directory=/app
command=/usr/local/bin/node build/index.js
